# src/provision/mb-provision.service
[Unit]
Description=MagicBridge WiFi Provisioning
After=NetworkManager.service network-pre.target
Before=magicbridge.service
# Runs on every boot and re-checks connectivity itself (see mb-provision.sh);
# no longer gated to "first boot only" so the setup hotspot can come back
# any time the Pi can't reach a known WiFi network (e.g. moved to a new
# location whose network was never saved).

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mb-provision.sh
RemainAfterExit=no
# Generous: provisioning BLOCKS on a human (join hotspot, type WiFi creds), and
# a wrong password now re-raises the hotspot for another attempt - 600s could
# kill it mid-retry and strand the unit again.
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
