# src/core/mb-mdns-alias.service
[Unit]
Description=MagicBridge optional mDNS alias (opt-in via config "mdns_alias"; off by default for stealth)
After=network-online.target avahi-daemon.service
Wants=network-online.target
Requires=avahi-daemon.service

[Service]
Type=simple
ExecStart=/usr/local/bin/mb-mdns-alias.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
